# templates/nimbus-supervisor.cfg
[program:nimbus]
command=/usr/local/storm/apache-storm-1.1.0/bin/storm nimbus
autorestart=true
stopsignal=KILL
stopasgroup=true
