# templates/zoo-supervisor.cfg
[program:zookeeper]
command=/usr/local/zookeeper/zookeeper-3.4.10/bin/zkServer.sh start-foreground
autorestart=true
stopsignal=KILL
stopasgroup=true
